# files/export-netcfg-thumbdrive.service
[Unit]
Description=Run script to export network configuration to a thumbdrive

[Service]
ExecStart=/home/pi/bin/export_network_config_thumb.sh
